# conf/supervisord.conf
[supervisord]
nodaemon=true
[program:nginx]
command=nginx
[program:php5fpm]
command=php5-fpm -c /etc/php5/fpm
[program:mysql]
command=mysqld
[program:mongodb]
command=mongod
